# scripts/Files/wireless.service
[Unit]
Description=Wireless profile service

[Service]
ExecStart=/etc/init.d/wireless

[Install]
WantedBy=multi-user.target