# scripts/Files/wired.service
[Unit]
Description=Wired profile service

[Service]
ExecStart=/etc/init.d/wired

[Install]
WantedBy=multi-user.target